# image-servers/djatoka/1.1.4/cfg/crontab
* * * * * root /usr/bin/curl -s 'http://localhost:8888/ingest?unattended=true' >>/var/log/cron.log 2>&1
#


